# deploy/live/etc/systemd/user/polyptic-agent.service.d/10-netboot-env.conf
# Drop-in for the user-scoped polyptic-agent.service (base unit written by `polyptic-agent setup`).
# Diskless netboot identity (POL-33): POLYPTIC_MACHINE_ID (from DMI, stable despite a random
# /etc/machine-id) + POLYPTIC_SERVER_URL / POLYPTIC_BOOTSTRAP_TOKEN (from the kernel cmdline) are written
# to this file by the SYSTEM oneshot polyptic-agent-env.service before greetd starts. These become REAL
# env vars, so they beat /etc/polyptic/agent.toml (applyConfigFileToEnv only fills UNSET vars). The
# leading '-' makes a missing file non-fatal (dev / non-netboot parity, the file only exists on a box
# that booted via the live image).
[Service]
EnvironmentFile=-/run/polyptic/agent.env
