# deploy/live/etc/systemd/system/systemd-networkd-wait-online.service.d/10-any-iface.conf
# Polyptic live image (POL-38): don't hold boot hostage to "ALL interfaces online".
#
# Stock wait-online blocks network-online.target until EVERY managed link is routable and only
# gives up after its own 120s timeout — on a kiosk with one live NIC (or a VM with an interface
# the image's netplan doesn't manage) that is a guaranteed ~2 minute stall on the cold-boot path,
# the opposite of the zero-click boot-to-content contract. One routable interface is all the agent
# needs, and it reconnects/retries forever anyway, so ANY iface + a short cap loses nothing.
[Service]
ExecStart=
ExecStart=/usr/lib/systemd/systemd-networkd-wait-online --any --timeout=15
