# deploy/agent/templates/polyptych-agent-firstboot.service
# polyptych-agent-firstboot.service — completes provisioning on first boot.  (systemd SYSTEM unit)
#   ->  installed to /usr/lib/systemd/system/polyptych-agent-firstboot.service
#
# Why this exists: a package's maintainer scripts cannot install Chromium during install, because
# the package-manager lock (dpkg/rpm) is held by the very transaction that's installing this package.
# So the package-manager half of `polyptych-agent setup` (deps incl. the real .deb/.rpm Chromium —
# NOT the snap, D27) runs here instead: once, on first boot, with the network up and no lock held.
# The unit disables itself after a successful run. Package-owned; the postinstall enables it.

[Unit]
Description=Polyptych first-boot provisioning (runs `polyptych-agent setup`)
Documentation=https://github.com/polyptych/polyptych/blob/main/docs/DEPLOY.md
Wants=network-online.target
After=network-online.target
# Skip once provisioning has completed (belt-and-braces with the self-disable below).
ConditionPathExists=!/var/lib/polyptych/.provisioned

[Service]
Type=oneshot
RemainAfterExit=no
# Full, idempotent provision: distro-detect (apt/dnf/pacman), install deps incl. the real Chromium,
# create the kiosk user, write greetd/sway/agent configs, enable greetd.
ExecStart=/usr/bin/polyptych-agent setup
# Record completion + stop re-running on subsequent boots.
ExecStartPost=/bin/sh -c 'mkdir -p /var/lib/polyptych && : > /var/lib/polyptych/.provisioned'
ExecStartPost=/bin/systemctl disable polyptych-agent-firstboot.service
# If setup fails (e.g. no network yet), retry on the next boot rather than giving up.
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
